# openwrt/init.d.webapp
#!/bin/sh /etc/rc.common
 
START=99
STOP=10
 
start() {        
	/bin/mount /root/en.wp.50k.articles.sqfs /root/mnt/articles/ -t squashfs
	/bin/mount /root/en.wp.50k.images.sqfs /root/mnt/images/ -t squashfs
	/bin/mount /root/en.wp.50k.index.sqfs /root/mnt/index/ -t squashfs
	/root/www/webappd start
}                 
 
stop() {          
	/root/www/webappd stop	
}	
